# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock    

[supervisord]
nodaemon=false   

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:nginx]
command=nginx
autorestart=true

[program:downloader]
command=python3 -c "import subprocess; import time; time.sleep(10); subprocess.check_call('supervisorctl start hound', shell=True)"
stdout_file=/dev/stdout
nodaemon=false

[program:hound]
command=python3 -c "import time; time.sleep(10000)"
autostart=false
